# install/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/onlyoffice/supervisord.log

[program:ASC.Sdk]
command=/usr/local/bin/node server.js
directory=/var/www/products/ASC.Sdk/sdk
autorestart=true
stdout_logfile=/var/log/onlyoffice/asc.sdk.out.log
stderr_logfile=/var/log/onlyoffice/asc.sdk.err.log
environment=NODE_ENV="production"

[program:ASC.Editors]
command=/usr/local/bin/node server.js
directory=/var/www/products/ASC.Editors/editor
autorestart=true
stdout_logfile=/var/log/onlyoffice/asc.editors.out.log
stderr_logfile=/var/log/onlyoffice/asc.editors.err.log
environment=NODE_ENV="production"

[program:ASC.Login]
command=/usr/local/bin/node server.js
directory=/var/www/products/ASC.Login/login
autorestart=true
stdout_logfile=/var/log/onlyoffice/asc.login.out.log
stderr_logfile=/var/log/onlyoffice/asc.login.err.log
environment=NODE_ENV="production"

[program:ASC.Socket.IO]
command=/usr/local/bin/node server.js
directory=/var/www/services/ASC.Socket.IO
autorestart=true
stdout_logfile=/var/log/onlyoffice/asc.socketio.out.log
stderr_logfile=/var/log/onlyoffice/asc.socketio.err.log
environment=NODE_ENV="production"

[program:ASC.SsoAuth]
command=/usr/local/bin/node app.js
directory=/var/www/services/ASC.SsoAuth
autorestart=true
stdout_logfile=/var/log/onlyoffice/asc.ssoauth.out.log
stderr_logfile=/var/log/onlyoffice/asc.ssoauth.err.log
environment=NODE_ENV="production"
