# install/docker/config/supervisor/dotnet_services.conf
[supervisord]
nodaemon=true
logfile=%(ENV_LOG_DIR)s/supervisord_dotnet.log
loglevel=warn

[program:ASC.Migration.Runner]
command=dotnet ASC.Migration.Runner.dll standalone=true
directory=%(ENV_BUILD_PATH)s/services/ASC.Migration.Runner/service
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=%(ENV_LOG_DIR)s/ASC.Migration.Runner.log
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_LOG_DIR)s/ASC.Migration.Runner.log
stderr_logfile_maxbytes=0

[eventlistener:start-after-migration]
command=python3 /usr/local/bin/start-after-migration.py
events=PROCESS_STATE_EXITED
autostart=true
autorestart=false
stderr_logfile=%(ENV_LOG_DIR)s/ASC.Migration.Runner.log
stderr_logfile_maxbytes=0

[program:ASC.AI]
command=dotnet ASC.AI.dll --urls=http://0.0.0.0:%(ENV_SERVICE_AI_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.AI core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/products/ASC.AI/server
autostart=false
autorestart=true
startsecs=25

[program:ASC.AI.Worker]
command=dotnet ASC.AI.Worker.dll --urls=http://0.0.0.0:%(ENV_SERVICE_AI_WORKER_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.AI.Worker core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server --core:eventBus:subscriptionClientName=asc_event_bus_ai_service_queue
directory=%(ENV_BUILD_PATH)s/products/ASC.AI/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.ApiSystem]
command=dotnet ASC.ApiSystem.dll --urls=http://0.0.0.0:%(ENV_SERVICE_API_SYSTEM_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.ApiSystem core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/services/ASC.ApiSystem/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.ClearEvents]
command=dotnet ASC.ClearEvents.dll --urls=http://0.0.0.0:%(ENV_SERVICE_CLEAR_EVENTS_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.ClearEvents core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/services/ASC.ClearEvents/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Data.Backup]
command=dotnet ASC.Data.Backup.dll --urls=http://0.0.0.0:%(ENV_SERVICE_BACKUP_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Data.Backup core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/services/ASC.Data.Backup/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Data.Backup.Worker]
command=dotnet ASC.Data.Backup.Worker.dll --urls=http://0.0.0.0:%(ENV_SERVICE_BACKUP_WORKER_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Data.Backup.Worker core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_backup_queue
directory=%(ENV_BUILD_PATH)s/services/ASC.Data.Backup.Worker/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Files]
command=dotnet ASC.Files.dll --urls=http://0.0.0.0:%(ENV_SERVICE_FILES_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Files core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/products/ASC.Files/server
autostart=false
autorestart=true
startsecs=25

[program:ASC.Files.Worker]
command=dotnet ASC.Files.Worker.dll --urls=http://0.0.0.0:%(ENV_SERVICE_FILES_WORKER_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Files.Worker core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_files_service_queue
directory=%(ENV_BUILD_PATH)s/products/ASC.Files/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Notify]
command=dotnet ASC.Notify.dll --urls=http://0.0.0.0:%(ENV_SERVICE_NOTIFY_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Notify core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_notify_queue
directory=%(ENV_BUILD_PATH)s/services/ASC.Notify/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.People]
command=dotnet ASC.People.dll --urls=http://0.0.0.0:%(ENV_SERVICE_PEOPLE_SERVER_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.People core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/products/ASC.People/server
autostart=false
autorestart=true
startsecs=25

[program:ASC.Studio.Notify]
command=dotnet ASC.Studio.Notify.dll --urls=http://0.0.0.0:%(ENV_SERVICE_STUDIO_NOTIFY_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Studio.Notify core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_studio_notify_queue
directory=%(ENV_BUILD_PATH)s/services/ASC.Studio.Notify/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.TelegramService]
command=dotnet ASC.TelegramService.dll --urls=http://0.0.0.0:%(ENV_SERVICE_TELEGRAM_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.TelegramService core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_telegram_queue
directory=%(ENV_BUILD_PATH)s/services/ASC.TelegramService/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Web.Api]
command=dotnet ASC.Web.Api.dll --urls=http://0.0.0.0:%(ENV_SERVICE_API_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.Api core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/services/ASC.Web.Api/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Web.HealthChecks.UI]
command=dotnet ASC.Web.HealthChecks.UI.dll --urls=http://0.0.0.0:%(ENV_SERVICE_HEALTHCHECKS_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.HealthChecks.UI core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server
directory=%(ENV_BUILD_PATH)s/services/ASC.Web.HealthChecks.UI/service
autostart=false
autorestart=true
startsecs=25

[program:ASC.Web.Studio]
command=dotnet ASC.Web.Studio.dll --urls=http://0.0.0.0:%(ENV_SERVICE_STUDIO_PORT)s --'$STORAGE_ROOT'=%(ENV_APP_STORAGE_ROOT)s --ENVIRONMENT=%(ENV_INSTALLATION_TYPE)s --pathToConf=%(ENV_PATH_TO_CONF)s --log:dir=%(ENV_LOG_DIR)s --log:name=ASC.Web.Studio core:products:folder=%(ENV_BUILD_PATH)s/products/ core:products:subfolder=server core:eventBus:subscriptionClientName=asc_event_bus_webstudio_queue
directory=%(ENV_BUILD_PATH)s/services/ASC.Web.Studio/service
autostart=false
autorestart=true
startsecs=25
