# run/macos/com.onlyoffice.webapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlyoffice.webapi</string>

    <key>WorkingDirectory</key>
    <string>/Users/alexeybannov/Documents/Works/docspace</string>
   
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/share/dotnet/dotnet</string>
        <string>server/web/ASC.Web.Api/bin/Debug/ASC.Web.Api.dll</string>
        <string>urls=http://0.0.0.0:5000</string>
        <string>$STORAGE_ROOT=/Users/alexeybannov/Documents/Works/docspace/Data</string>
        <string>log:dir=/Users/alexeybannov/Documents/Works/docspace/Logs</string>
        <string>log:name=web.api</string>
        <string>pathToConf=/Users/alexeybannov/Documents/Works/docspace/buildtools/config</string>
        <string>core:products:folder=/Users/alexeybannov/Documents/Works/docspace/server/products</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <true/>

    <!--<key>StandardOutPath</key>
    <string>/Users/alexeybannov/Documents/Works/docspace/Logs/webapi.stdout.log</string>-->

    <!--<key>StandardErrorPath</key>
    <string>/Users/alexeybannov/Documents/Works/docspace/Logs/webapi.stderr.log</string>-->
</dict>
</plist>
