# run/macos/com.onlyoffice.client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlyoffice.client</string>

    <key>WorkingDirectory</key>
    <string>${DOCSPACE_ROOT}/client/packages/client</string>
   
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>server.js</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <!--<key>StandardErrorPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.client.err</string>-->

    <!--<key>StandardOutPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.client.log</string>-->

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
