# run/macos/com.onlyoffice.apisystem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlyoffice.apisystem</string>

    <key>WorkingDirectory</key>
    <string>${DOCSPACE_ROOT}</string>
   
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/share/dotnet/dotnet</string>
        <string>server/common/services/ASC.ApiSystem/bin/Debug/ASC.ApiSystem.dll</string>
        <string>urls=http://0.0.0.0:5010</string>
        <string>$STORAGE_ROOT=${DOCSPACE_ROOT}/Data</string>
        <string>log:dir=${DOCSPACE_ROOT}/Logs</string>
        <string>log:name=apisystem</string>
        <string>pathToConf=${DOCSPACE_ROOT}/buildtools/config</string>
        <string>core:products:folder=${DOCSPACE_ROOT}/server/products</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.apisystem.err</string>

    <key>StandardOutPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.apisystem.log</string>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
