# run/macos/com.onlyoffice.socketio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlyoffice.socketio</string>

    <key>WorkingDirectory</key>
    <string>${DOCSPACE_ROOT}</string>
   
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>server/common/ASC.Socket.IO/server.js</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardErrorPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.socketio.err</string>

    <key>StandardOutPath</key>
    <string>${DOCSPACE_ROOT}/Logs/com.onlyoffice.socketio.log</string>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
